# systemd/.config/systemd/user/backup-ocean.service
[Unit]
Description=Backup by using borg

[Service]
Type=simple
ExecStart=/bin/bash -c %h/bin/backup-ocean.sh

[Install]
WantedBy=multi-user.target
