# smart-audio/99-smart-audio.rules
# /etc/udev/rules.d/99-smart-audio.rules

# Rule for Bluetooth device connect/disconnect
# Uses bluetooth subsystem with DEVTYPE=link which fires for any BT device connection.
# The script handles prioritization (it checks for bluez_output sinks).
ACTION=="add", SUBSYSTEM=="bluetooth", ENV{DEVTYPE}=="link", RUN+="/usr/local/bin/smart-audio-switcher.sh --device-connect bluetooth"
ACTION=="remove", SUBSYSTEM=="bluetooth", ENV{DEVTYPE}=="link", RUN+="/usr/local/bin/smart-audio-switcher.sh --device-disconnect"

# Rule for any display plug/unplug event on the graphics card
# This now calls the main script directly to perform the check.
ACTION=="change", SUBSYSTEM=="drm", KERNEL=="card*", RUN+="/usr/local/bin/smart-audio-switcher.sh --hdmi-check"
