# executables/.config/systemd/user/backup.timer
[Unit]
Description=Backup timer

[Timer]
OnBootSec=15min
Persistent=true
OnCalendar=daily
Unit=backup.service

[Install]
WantedBy=timers.target
